# resources/puppetlabs/lein-ezbake/template/pe/ext/redhat/ezbake.service.erb
#
# Local settings can be configured without being overwritten by package upgrades, for example
# if you want to increase <%= EZBake::Config[:project] %> open-files-limit to 10000,
# you need to increase systemd's LimitNOFILE setting, so create a file named
# "/etc/systemd/system/<%= EZBake::Config[:project] %>.service.d/limits.conf" containing:
#	[Service]
#	LimitNOFILE=10000
# You can confirm it worked by running systemctl daemon-reload
# then running systemctl show <%= EZBake::Config[:project] %> | grep LimitNOFILE
#
[Unit]
Description=<%= EZBake::Config[:project] %> Service
After=syslog.target network.target <%= EZBake::Config[:start_after].map {|dep| "#{dep}.service" }.join(" ") %>

[Service]
Type=simple
EnvironmentFile=/etc/sysconfig/<%= EZBake::Config[:project] %>
User=<%= EZBake::Config[:user] %>
TimeoutStartSec=<%= EZBake::Config[:start_timeout] %>
TimeoutStopSec=60
Restart=on-failure
StartLimitBurst=5
#set default privileges to -rw-r-----
UMask=027
<% if EZBake::Config[:open_file_limit] -%>
LimitNOFILE=<%= EZBake::Config[:open_file_limit] %>
<% end -%>

<% unless EZBake::Config[:redhat][:pre_start_action].empty? -%>
PermissionsStartOnly=true
<% EZBake::Config[:redhat][:pre_start_action].each do |action| -%>
ExecStartPre=<%= action %>
<% end -%>
<% end -%>

ExecStart=/opt/puppetlabs/server/bin/java $JAVA_ARGS \
          '-XX:OnOutOfMemoryError=kill -9 %%p' \
          -Djava.security.egd=/dev/urandom \
          -cp "${INSTALL_DIR}/<%= EZBake::Config[:uberjar_name] %>" clojure.main \
          -m <%= EZBake::Config[:main_namespace] %> \
          --config "${CONFIG}" \
          -b "${BOOTSTRAP_CONFIG}" $@

KillMode=process

ExecStartPost=/bin/bash "${INSTALL_DIR}/ezbake-functions.sh" wait_for_app

<% EZBake::Config[:redhat][:post_start_action].each do |action| -%>
ExecStartPost=-<%= action %>
<% end -%>
SuccessExitStatus=143

StandardOutput=syslog

[Install]
WantedBy=multi-user.target
