# resources/puppetlabs/lein-ezbake/template/pe/ext/redhat/ezbake.service.erb
[Unit]
Description=<%= EZBake::Config[:project] %> Service
After=syslog.target network.target <%= EZBake::Config[:start_after].map {|dep| "#{dep}.service" }.join(" ") %>

[Service]
Type=forking
EnvironmentFile=/etc/sysconfig/<%= EZBake::Config[:project] %>
User=<%= EZBake::Config[:user] %>
TimeoutStartSec=<%= EZBake::Config[:start_timeout] %>
TimeoutStopSec=60
Restart=on-failure
StartLimitBurst=5
PIDFile=/var/run/puppetlabs/<%= EZBake::Config[:real_name] %>/<%= EZBake::Config[:real_name] %>
#set default privileges to -rw-r-----
UMask=027
<% if EZBake::Config[:open_file_limit] -%>
LimitNOFILE=<%= EZBake::Config[:open_file_limit] %>
<% end -%>

<% unless EZBake::Config[:redhat][:pre_start_action].empty? -%>
PermissionsStartOnly=true
<% EZBake::Config[:redhat][:pre_start_action].each do |action| -%>
ExecStartPre=<%= action %>
<% end -%>
<% end -%>

ExecReload=/opt/puppetlabs/server/apps/<%= EZBake::Config[:real_name] %>/bin/<%= EZBake::Config[:real_name] %> reload
ExecStart=/opt/puppetlabs/server/apps/<%= EZBake::Config[:real_name] %>/bin/<%= EZBake::Config[:real_name] %> start

KillMode=process

<% EZBake::Config[:redhat][:post_start_action].each do |action| -%>
ExecStartPost=-<%= action %>
<% end -%>
SuccessExitStatus=143

StandardOutput=syslog

[Install]
WantedBy=multi-user.target
