# template/foss/ext/redhat/ezbake.service.erb
[Unit]
Description=<%= EZBake::Config[:project] %> Service
After=syslog.target network.target

[Service]
Type=simple
EnvironmentFile=/etc/sysconfig/<%= EZBake::Config[:project] %>
User=<%= EZBake::Config[:user] %>

ExecStart=/usr/bin/java $JAVA_ARGS \
          -cp "${INSTALL_DIR}/<%= EZBake::Config[:uberjar_name] %>" clojure.main \
          -m puppetlabs.trapperkeeper.main --config "${CONFIG}" \
          -b "${BOOTSTRAP_CONFIG}" $@

ExecStop=/bin/kill $MAINPID

StandardOutput=syslog

[Install]
WantedBy=multi-user.target
