# template/foss/ext/redhat/ezbake.service.erb
[Unit]
Description=<%= EZBake::Config[:project] %> Service
After=syslog.target network.target

[Service]
Type=simple
EnvironmentFile=/etc/sysconfig/<%= EZBake::Config[:project] %>
User=<%= EZBake::Config[:user] %>
TimeoutStopSec=60

ExecStart=/usr/bin/java $JAVA_ARGS \
          -XX:OnOutOfMemoryError=kill\ -9\ %p \
          -XX:+HeapDumpOnOutOfMemoryError \
          -XX:HeapDumpPath=/var/log/<%= EZBake::Config[:project] %> \
          -cp "${INSTALL_DIR}/<%= EZBake::Config[:uberjar_name] %>" clojure.main \
          -m puppetlabs.trapperkeeper.main --config "${CONFIG}" \
          -b "${BOOTSTRAP_CONFIG}" $@

KillMode=process

ExecStartPost=/bin/bash /usr/share/puppetserver/ezbake-functions.sh wait_for_app

<% EZBake::Config[:redhat][:post_start_action].each do |action| -%>
ExecStartPost=-<%= action %>
<% end -%>
SuccessExitStatus=143

StandardOutput=syslog

[Install]
WantedBy=multi-user.target
