# resources/puppetlabs/lein-ezbake/template/pe/ext/debian/prerm.erb
#!/bin/sh

if [ "$1" = failed-upgrade ]; then
    invoke-rc.d <%= EZBake::Config[:project] -%> stop || :
fi
