# resources/puppetlabs/lein-ezbake/template/foss/ext/debian/prerm.erb
#!/bin/sh

if [ "$1" = "failed-upgrade" ] || [ "$1" = "remove" ] ; then
  if [ -d /run/systemd/system ]; then
    # Using systemd
    systemctl --no-reload disable <%= EZBake::Config[:project] -%>.service > /dev/null 2>&1 || :
    systemctl stop <%= EZBake::Config[:project] -%>.service > /dev/null 2>&1 || :
  else
    # Using sysv
    invoke-rc.d <%= EZBake::Config[:project] -%> stop || :
  fi
fi
