# resources/puppetlabs/lein-ezbake/template/global/ext/debian/ezbake.service.erb
#
# Local settings can be configured without being overwritten by package upgrades, for example
# if you want to increase <%= EZBake::Config[:project] %> open-files-limit to 10000,
# you need to increase systemd's LimitNOFILE setting, so create a file named
# "/etc/systemd/system/<%= EZBake::Config[:project] %>.service.d/limits.conf" containing:
#	[Service]
#	LimitNOFILE=10000
# You can confirm it worked by running systemctl daemon-reload
# then running systemctl show <%= EZBake::Config[:project] %> | grep LimitNOFILE
#
[Unit]
Description=<%= EZBake::Config[:project] %> Service
<% if EZBake::Config[:start_before].any? -%>
Before=<%= EZBake::Config[:start_before].map {|dep| "#{dep}.service" }.join(" ") %>
<% end -%>
After=syslog.target network.target nss-lookup.target <%= EZBake::Config[:start_after].map {|dep| "#{dep}.service" }.join(" ") %>
<% if EZBake::Config[:start_after].any? -%>
Wants=<%= EZBake::Config[:start_after].map {|dep| "#{dep}.service" }.join(" ") %>
<% end -%>

[Service]
Type=forking
EnvironmentFile=/etc/default/<%= EZBake::Config[:project] %>
User=<%= EZBake::Config[:user] %>
TimeoutStartSec=<%= EZBake::Config[:start_timeout] %>
TimeoutStopSec=<%= EZBake::Config[:stop_timeout] %>
Restart=on-failure
StartLimitBurst=5
PIDFile=/run/puppetlabs/<%= EZBake::Config[:real_name] %>/<%= EZBake::Config[:real_name] %>.pid

# https://tickets.puppetlabs.com/browse/EZ-129
# Prior to systemd v228, TasksMax was unset by default, and unlimited. Starting in 228 a default of '512'
# was implemented. This is low enough to cause problems for certain applications. In systemd 231, the
# default was changed to be 15% of the default kernel limit. This explicitly sets TasksMax to 4915,
# which should match the default in systemd 231 and later.
# See https://github.com/systemd/systemd/issues/3211#issuecomment-233676333
TasksMax=4915

#set default privileges to -rw-r-----
UMask=027
<% if EZBake::Config[:open_file_limit] -%>
LimitNOFILE=<%= EZBake::Config[:open_file_limit] %>
<% end -%>

<% unless EZBake::Config[:debian][:pre_start_action].empty? -%>
PermissionsStartOnly=true
<% EZBake::Config[:debian][:pre_start_action].each do |action| -%>
ExecStartPre=<%= action %>
<% end -%>
<% end -%>

ExecReload=/opt/puppetlabs/server/apps/<%= EZBake::Config[:real_name] %>/bin/<%= EZBake::Config[:real_name] %> reload
ExecStart=/opt/puppetlabs/server/apps/<%= EZBake::Config[:real_name] %>/bin/<%= EZBake::Config[:real_name] %> start
ExecStop=/opt/puppetlabs/server/apps/<%= EZBake::Config[:real_name] %>/bin/<%= EZBake::Config[:real_name] %> stop

KillMode=process

<% EZBake::Config[:debian][:post_start_action].each do |action| -%>
ExecStartPost=-<%= action %>
<% end -%>
SuccessExitStatus=143

[Install]
WantedBy=multi-user.target
