# resources/puppetlabs/lein-ezbake/template/global/ext/debian/postinst.erb
#!/bin/sh

/opt/puppetlabs/server/apps/<%= EZBake::Config[:real_name] -%>/scripts/install.sh postinst_deb

# On upgrade, we should restart the service if it's running
if [ $1 = 'configure' -a -n $2 ] ; then
  if [ -d '/run/systemd/system' ] ; then
    # Using systemd
    systemctl daemon-reload >/dev/null 2>&1 || :
    systemctl try-restart <%= EZBake::Config[:project] %>.service ||:
  else
    # Using sysv
    invoke-rc.d <%= EZBake::Config[:project] %> try-restart || :
  fi
fi

# Run trigger scripts on install if defined
if [ $1 = triggered ] && [ -z $2 ]; then
  <% EZBake::Config[:debian][:interested_install_triggers].each do |_, scripts| -%>
    <%= scripts.join("\n") %>
  <% end -%>
  : # in case there are no install triggers
fi

# Run trigger scripts on upgrade if defined
if [ $1 = triggered ] && [ -n $2 ]; then
  <% EZBake::Config[:debian][:interested_upgrade_triggers].each do |_, scripts| -%>
    <%= scripts.join("\n") %>
  <% end -%>
  : # in case there are no upgrade triggers
fi
