# template/ext/debian/prerm.erb
#!/bin/sh

set -e
if [ -x "/etc/init.d/<%= EZBake::Config[:project] %>" ]; then
    if [ -x "`which invoke-rc.d 2>/dev/null`" ]; then
        invoke-rc.d <%= EZBake::Config[:project] %> stop || exit $?
    else
        /etc/init.d/<%= EZBake::Config[:project] %> stop || exit $?
    fi
fi

